# install/doorbell.service
[Unit]
Description=Doorbell
After=network.target

[Service]
ExecStart=/usr/lib/nerdhome/.env/bin/python3 /usr/lib/nerdhome/doorbell.py
Restart=always
RestartSec=5
KillMode=process
Type=simple

[Install]
WantedBy=multi-user.target
